# .config/autostart/yakuake.desktop
[Desktop Entry]
Name=Yakuake
Exec=/usr/bin/yakuake
Type=Application
#Terminal=False
